# scripts/base/systemd/ryzenadj-performance.timer
[Unit]
Description="Run ryzenadj-performance.service after 3sec waitting"

[Timer]
OnActiveSec=3
Unit=ryzenadj-performance.service
