# scripts/base/systemd/ryzenadj-powersave.service
[Unit]
Description="Change amd powermode to powersave"

[Service]
ExecStart=/usr/bin/ryzenadj --power-saving
