# scripts/base/systemd/ryzenadj-performance.service
[Unit]
Description="Change amd powermode to performance"

[Service]
ExecStart=/usr/bin/ryzenadj --max-performance
